# templates/startup/server.conf.j2
# Run prometheus service

start on startup

script
   {{prometheus_base_dir}}/server/prometheus -config.file={{prometheus_config_path}}/prometheus.yml
end script
